# systemd/android-automount@.service
# Android MTP mount service

[Service]
Type=oneshot
ExecStart=/usr/bin/android-mount %I
ExecStop=/usr/bin/android-umount %I
RemainAfterExit=yes
