# systemd/media-automount@.service
[Service]
Type=forking
ExecStart=/usr/bin/media-mount %I
ExecStop=/usr/bin/media-umount
