# systemd/media-automount@.service
[Service]
Type=oneshot
ExecStart=/usr/bin/media-mount %I
ExecStop=/usr/bin/media-umount
RemainAfterExit=yes
